# kittycam/kittycam.service
[Unit]
Description=Runs the Kitty Cam NodeJS service
After=network.target
Requires=raspivid.service

[Service]
EnvironmentFile=/usr/local/KittyCam/kittycam/systemd.conf
WorkingDirectory=/usr/local/KittyCam/kittycam
ExecStart=/usr/bin/node kittyCam.js ${CONFIG} ${IMAGE}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
